# backup/luci-app-qosv4/root/etc/init.d/qosv4
#!/bin/sh /etc/rc.common

START=60

start(){
	/usr/bin/qosv4 start
}

stop(){
	/usr/bin/qosv4 stop
}
restart(){
	/etc/init.d/qosv4 stop
	/etc/init.d/qosv4 start
}



